# rootfs-overlay/etc/inittab
::sysinit:/sbin/openrc sysinit
::sysinit:/sbin/openrc boot
::wait:/sbin/openrc default
tty1::respawn:/sbin/agetty --noclear 38400 tty1 linux
tty2::respawn:/sbin/agetty 38400 tty2 linux
tty3::respawn:/sbin/agetty 38400 tty3 linux
::ctrlaltdel:/sbin/reboot
::shutdown:/sbin/openrc shutdown

